# ansible/roles/vault_agent/templates/containers/vault-agent.container
# Copyright (C) - LOW-LAYER
# Contact : contact@low-layer.com

# ============================================================================
# Systemd unit file for running HashiCorp Vault Agent as a quadlet container.
# This setup allows Vault Agent to authenticate with Vault using AppRole
# secrets, render templates, and manage certificates securely.
# ============================================================================

[Unit]
Description=Vault Agent container
Wants=network-online.target
After=network-online.target

[Container]
Image=docker.io/hashicorp/vault:1.15
ContainerName=vault-agent

# ============================================================================
# Command Execution
# ============================================================================
Exec=agent -config=/vault/conf/agent.hcl

# ============================================================================
# Network configuration
# ============================================================================
Network=host

# ============================================================================
# Capabilities
# ============================================================================
# Prevent memory from being swapped
AddCapability=IPC_LOCK

# Allow file capabilities to be set
AddCapability=CAP_SETFCAP

# Required for managing file ownership
AddCapability=CHOWN

# ============================================================================
# Volumes
# ============================================================================
Volume={{ system_ssl_certs_path }}:{{ system_ssl_certs_path }}:ro
Volume={{ system_pki_ca_trust_path }}:{{ system_pki_ca_trust_path }}:ro
Volume={{ vault_agent_config_path }}:/vault/conf:ro
Volume={{ vault_agent_templates_path }}:/vault/templates:ro
Volume={{ vault_agent_certs_path }}:/certs:z
# Volume=/etc/localtime:/etc/localtime:ro

# ============================================================================
# Podman secrets mounted in /run/secrets
# ============================================================================
Secret=vault-role-id,type=mount,target=/run/secrets/vault_role_id
Secret=vault-secret-id,type=mount,target=/run/secrets/vault_secret_id

# ============================================================================
# Security configuration
# ============================================================================
PodmanArgs=--user=0:0 --entrypoint=vault
NoNewPrivileges=false

# ============================================================================
# Resource management
# ============================================================================
Memory=1G

[Service]
Restart=always
RestartSec=10
TimeoutStartSec=900
TimeoutStopSec=70

[Install]
WantedBy=multi-user.target
